# config/vps-monitor.service
[Unit]
Description=VPS库存监控机器人
Documentation=https://github.com/kure29/VPSMonitorBot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=vpsmonitor
Group=vpsmonitor
WorkingDirectory=/opt/vps-monitor
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/vps-monitor/venv/bin/python3 /opt/vps-monitor/monitor.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-monitor

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vps-monitor/data /opt/vps-monitor/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# 资源限制
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=512M
CPUQuota=50%

# 环境变量文件
EnvironmentFile=-/opt/vps-monitor/vps-monitor.env

[Install]
WantedBy=multi-user.target

# 安装说明：
# 1. 将此文件复制到 /etc/systemd/system/vps-monitor.service
# 2. 创建用户和目录：
#    sudo useradd -r -s /bin/false vpsmonitor
#    sudo mkdir -p /opt/vps-monitor
#    sudo chown vpsmonitor:vpsmonitor /opt/vps-monitor
# 3. 将项目文件复制到 /opt/vps-monitor/
# 4. 设置权限：
#    sudo chown -R vpsmonitor:vpsmonitor /opt/vps-monitor
# 5. 启用并启动服务：
#    sudo systemctl daemon-reload
#    sudo systemctl enable vps-monitor.service
#    sudo systemctl start vps-monitor.service
# 6. 查看状态：
#    sudo systemctl status vps-monitor.service
#    sudo journalctl -u vps-monitor.service -f
